# build/Python/src/nfc/trusty_debian/nfc.conf
description  "Near Field Communication server"
author    "ProteinSimple"

start on runlevel [2345]
stop on (runlevel [016] or starting cartridge_editor)

console log
respawn
setuid cbuser
setgid cbuser
script
    chdir /home/cbuser/cellbio/nfc
    rm -f /home/cbuser/cellbio/nfc/crossbar/node.pid
    exec /usr/share/nfc/bin/crossbar start --cbdir /home/cbuser/cellbio/nfc/crossbar
end script

